# go2rtc/com.xiaomi-token-watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.xiaomi-token-watch</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/etmacmini/代码/smartHome2/scripts/xiaomi-token-watch.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/etmacmini/代码/smartHome2/logs/xiaomi-token-watch.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/etmacmini/代码/smartHome2/logs/xiaomi-token-watch.error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/etmacmini/代码/smartHome2</string>
</dict>
</plist>
